# scripts/foli-harvester.service
[Unit]
Description=Foli SIRI data harvester
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/omat/foli-data-collector
ExecStart=%h/omat/foli-data-collector/.venv/bin/foli-harvester collect
Restart=always
RestartSec=15

[Install]
WantedBy=default.target
